# scripts/ssetop-server.service
[Unit]
Description=ssetop server

[Service]
ExecStart=/usr/bin/ssetop-server

[Install]
WantedBy=multi-user.target
